# tools/systemd/php8.4-fpm.service.d/pantedu-sandbox.conf
# Phase 25.K.5 v3 — PHP-FPM systemd sandboxing (Debian Trixie + PHP 8.4).
# Validato via incremental testing — config minimal funzionante.
#
# Test:
#   systemctl restart php8.4-fpm  # active atteso
#   curl -sS https://yourapp/     # 200/302 atteso
#
# Rimosse direttive incompatibili rispetto v1 originale:
# - ProtectProc=invisible + ProcSubset=pid (master vedeva workers via /proc → SIGSYS)
# - MemoryDenyWriteExecute=true (PHP runtime bytecode gen)
# - SystemCallFilter=@system-service (filter restrittivo → SIGSYS status=31)
# - CapabilityBoundingSet=... (alcune capabilities necessarie a master fork)
# - ProtectSystem=strict (rompe accesso /etc/php/conf.d/* runtime; full OK)
# - RestrictNamespaces (PHP usa namespaces interni)
# - RestrictAddressFamilies (richiesta da CrowdSec curl, HIBP, threat-intel)
#
# Rollback:
#   sudo rm /etc/systemd/system/php8.4-fpm.service.d/pantedu-sandbox.conf
#   sudo systemctl daemon-reload && sudo systemctl restart php8.4-fpm

[Service]
# ── Filesystem protection ──
PrivateTmp=true
ProtectHome=true
ProtectSystem=full              # /usr, /boot read-only (NOT strict)

# ── Privileges ──
NoNewPrivileges=true              # no setuid escalation
RestrictSUIDSGID=true             # no creation suid/sgid files

# ── Kernel protection ──
ProtectKernelTunables=true        # /proc/sys read-only
ProtectKernelLogs=true            # /proc/kmsg, dmesg negato
ProtectControlGroups=true         # /sys/fs/cgroup read-only

# ── No clock/hostname tampering ──
ProtectClock=true
ProtectHostname=true

# ── Device protection ──
PrivateDevices=true               # no /dev/sd*, /dev/mem (solo basic /dev/null, /dev/random)

# ── No realtime / personality ──
RestrictRealtime=true             # no realtime scheduling
LockPersonality=true              # no personality() syscall (anti-architecture-spoof)

# ── Resource limits ──
LimitNOFILE=8192
