# tools/systemd/php8.4-fpm.service.d/pantedu-auto-deploy.conf
# Phase 25.R.21 — PHP-FPM drop-in per auto-deploy webhook (Opzione 4).
#
# Aggiunge /var/lib/pantedu-deploy alla whitelist ReadWritePaths.
# tools/webhook/github.php scrive il flag file in questa dir dopo HMAC
# verify. systemd Path unit (pantedu-deploy.path) watcha il file e
# lancia il deploy script come root.
#
# Drop-in separato dal pantedu-sandbox.conf esistente per separation
# of concerns: se rollback auto-deploy serve, basta rimuovere questo
# singolo file senza toccare la sandbox base.
#
# Install: copy in /etc/systemd/system/php8.4-fpm.service.d/ +
# systemctl daemon-reload + restart php8.4-fpm.
# Auto-installato da tools/webhook/install_auto_deploy.sh.

[Service]
ReadWritePaths=/var/lib/pantedu-deploy
